# home_jotah/denis_unified_v1/deploy/denis-cp-daemon.service
[Unit]
Description=Denis Control Plane Daemon - Human-in-the-loop approval
After=network.target

[Service]
Type=simple
WorkingDirectory=/media/jotah/SSD_denis/home_jotah/denis_unified_v1
Environment="PYTHONPATH=/media/jotah/SSD_denis/home_jotah/denis_unified_v1"
Environment="DISPLAY=:0"
ExecStart=/usr/bin/python3 /media/jotah/SSD_denis/home_jotah/denis_unified_v1/scripts/denis_cp_daemon.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
